# scripts/install/systemd/dev-studio-watcher-reload.service.tmpl
[Unit]
Description=Restart all dev-studio watchers when agent-watch.sh changes
Documentation=https://github.com/{{GITHUB_OWNER}}/{{GITHUB_REPO}}/blob/main/docs/decisions/ADR-0006-watcher-resilience.md
# Debounce: a git pull may touch the file briefly twice; cap reload churn
StartLimitIntervalSec=30
StartLimitBurst=2

[Service]
Type=oneshot
# Sleep absorbs sub-second multi-write bursts from editor/git so we issue one restart.
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/systemctl --user restart \
    dev-studio-watcher@orchestrator.service \
    dev-studio-watcher@product-manager.service \
    dev-studio-watcher@architect.service \
    dev-studio-watcher@developer.service \
    dev-studio-watcher@tester.service
# Log the event so post-mortem is easy
ExecStartPost=/bin/sh -c 'echo "[watcher-reload] $(date -u +%%Y-%%m-%%dT%%H:%%M:%%SZ) restarted all 5 watchers (agent-watch.sh changed)" >> /var/log/dev-studio/watcher-reload.log'
