# scripts/install/systemd/dev-studio-watcher-reload.path.tmpl
[Unit]
Description=Watch agent-watch.sh for changes and reload all watchers
Documentation=https://github.com/{{GITHUB_OWNER}}/{{GITHUB_REPO}}/blob/main/docs/decisions/ADR-0006-watcher-resilience.md

[Path]
# Fires on IN_CLOSE_WRITE + IN_MOVED_TO (atomic replace via git pull, editor save)
PathChanged={{REPO_ROOT}}/scripts/agent-watch.sh
# Trigger the reload service; the service handles the actual restart fan-out
Unit=dev-studio-watcher-reload.service

[Install]
WantedBy=default.target
