# systemd/dev-studio-health.service
[Unit]
Description=Dev Studio periodic health probe
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=atilcan
Group=atilcan
WorkingDirectory=/opt/dev-studio/atilprojects
EnvironmentFile=/home/atilcan/.dev-studio-env
ExecStart=/opt/dev-studio/atilprojects/scripts/health-check.sh
StandardOutput=journal
StandardError=journal
# Don't restart on failure — timer handles cadence
Restart=no

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/dev-studio
